# deploy/uav-phm.service
[Unit]
Description=UAV PHM Monitor — real-time fault detection and health monitoring
After=multi-user.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/drone_monitor
ExecStart=/home/pi/drone_monitor/venv/bin/python /home/pi/drone_monitor/monitor.py --hardware
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment=PHM_PORT=5000

[Install]
WantedBy=multi-user.target
